# timers/zfs-auto-daily.timer
[Unit]
Description=Daily ZFS snapshots at 1:20AM.

[Timer]
OnCalendar=*-*-* 01:20:00
