# timers/zfs-auto-snapshot-monthly.service
Description=ZFS monthly snapshot service

[Service]
ExecStart=/usr/local/sbin/zfs-auto-snapshot --quiet --syslog --label=monthly --keep=12 //
