# timers/zfs-auto-snapshot-monthly.service
[Unit]
Description=ZFS monthly snapshot service

[Service]
ExecStart=zfs-auto-snapshot  --quiet --syslog --label=monthly --keep=12 //
