# timers/zfs-auto-snapshot-weekly.timer
[Unit]
Description=ZFS weekly snapshot timer

[Timer]
OnCalendar=weekly
Persistent=true
RandomizedDelaySec=600
