# timers/zfs-auto-snapshot-hourly.service
[Unit]
Description=ZFS hourly snapshot service

[Service]
ExecStart=/usr/local/sbin/zfs-auto-snapshot --quiet --syslog --label=hourly --keep=24 //
