# timers/zfs-auto-snapshot-weekly.service
[Unit]
Description=ZFS weekly snapshot service

[Service]
ExecStart=/usr/local/sbin/zfs-auto-snapshot --quiet --syslog --label=weekly --keep=8 //
