# timers/zfs-auto-snapshot-daily.service
[Unit]
Description=ZFS daily snapshot service

[Service]
ExecStart=/usr/local/sbin/zfs-auto-snapshot --quiet --syslog --label=daily --keep=31 //
