# etc/zfs-auto-snapshot.cron.daily
#!/bin/sh
exec zfs-auto-snapshot --quiet --syslog --label=daily --keep=31 //
