# etc/zfs-auto-snapshot.cron.weekly
#!/bin/sh

# Only call zfs-auto-snapshot if it's available
exec which zfs-auto-snapshot > /dev/null && \
     zfs-auto-snapshot --quiet --syslog --label=weekly --keep=8 //
