# etc/zfs-auto-snapshot.cron.daily
#!/bin/sh
PATH="/usr/bin:/bin:/usr/sbin:/sbin:/usr/local/sbin"
exec zfs-auto-snapshot --quiet --syslog --label=daily --keep=31 //
