# etc/zfs-auto-snapshot.cron.hourly
#!/bin/sh

PATH="/usr/bin:/bin:/usr/local/sbin:/usr/sbin:/sbin"

# Only call zfs-auto-snapshot if it's available
which zfs-auto-snapshot > /dev/null || exit 0

exec zfs-auto-snapshot --quiet --syslog --label=hourly --keep=24 //
