# etc/zfs-auto-snapshot.cron.monthly
#!/bin/sh

# Only call zfs-auto-snapshot if it's available
which zfs-auto-snapshot > /dev/null && \
     zfs-auto-snapshot --quiet --syslog --label=monthly --keep=12 //
