# etc/zfs-auto-snapshot.cron.daily
#!/bin/sh
exec /usr/local/sbin/zfs-auto-snapshot --quiet --syslog --label=daily --keep=31 //
