# Library/org.zfs.snapshot.hourly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>org.zfs.snapshot.hourly</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/sbin/zfs-auto-snapshot.sh</string>
		<string>--syslog</string>
		<string>--keep=24</string>
		<string>--label=hourly</string>
		<string>//</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StartInterval</key>
	<integer>3600</integer>
</dict>
</plist>
