# src/samples/auto-snapshot-space-timf,frequent.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<service_bundle type='manifest' name='space-timf,frequent'>

<!-- A service instance that takes snapshots every 10 minutes, keeping 60 snapshots
  into the past. All child filesystems of space/timf are also snapshotted. -->

<service
	name='system/filesystem/zfs/auto-snapshot'
	type='service'
	version='0.10'>
	<create_default_instance enabled='false' />

	<instance name='space-timf,frequent' enabled='false' >

        <exec_method
		type='method'
		name='start'
		exec='/lib/svc/method/zfs-auto-snapshot start'
		timeout_seconds='0' />

	<exec_method
		type='method'
		name='stop'
		exec='/lib/svc/method/zfs-auto-snapshot stop'
		timeout_seconds='0' />

        <property_group name='startd' type='framework'>
        	<propval name='duration' type='astring' value='transient' />
        </property_group>

        <!-- properties for zfs automatic snapshots -->
	<property_group name="zfs" type="application">

	  <propval name="fs-name" type="astring" value="space/timf" 
		   override="true"/>
	  <propval name="interval" type="astring" value="minutes"
		   override="true"/>
	  <propval name="period" type="astring" value="10"
		   override="true"/>
	  <propval name="offset" type="astring" value="0"
		   override="true"/>
	  <propval name="keep" type="astring" value="60"
		   override="true"/>
	  <propval name="snapshot-children" type="boolean" value="true"
		   override="true"/>

	  <propval name="backup" type="astring" value="none"
		   override="true"/>
	  <propval name="backup-save-cmd" type="astring" value=""
		   override="true"/>
	  <propval name="backup-lock" type="astring" value="unlocked"
		   override="true"/>

	  <propval name="label" type="astring" value="frequent"
		   override="true"/>
          <propval name="verbose" type="boolean" value="false"
                   override="true"/>
	  <propval name="avoidscrub" type="boolean" value="true"
		   override="true"/>

	</property_group>

	</instance>

	<stability value='Unstable' />
</service>
</service_bundle>
