# example-timers/zfs-auto-daily.timer
[Unit]
Description=Daily ZFS snapshots at 1:15AM.

[Timer]
OnCalendar=*-*-* 01:15:00
 
[Install]
WantedBy=timers.target
