# example-timers/zfs-auto-weekly.timer
[Unit]
Description=Weekly ZFS snapshots.

[Timer]
OnCalendar=weekly
 
[Install]
WantedBy=timers.target
