# example-timers/zfs-auto-daily.timer
[Unit]
Description=Daily ZFS snapshots.

[Timer]
OnCalendar=daily
 
[Install]
WantedBy=timers.target
